# crontab
0 */6 * * * sudo systemctl restart tt-bot
