# crontab
0 */12 * * * sudo systemctl restart tt-bot
