# systemd/tt-bot-restart.service
[Unit]
Description=restart TT Bot service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart tt-bot
